# systemd/blit-share@.service
[Unit]
Description=blit share %i
After=network.target

[Service]
Type=notify
EnvironmentFile=/etc/blit/share-%i.env
ExecStart=blit share --quiet
Restart=on-failure

[Install]
WantedBy=multi-user.target
