# systemd/blit-server.service
[Unit]
Description=blit terminal multiplexer
Requires=blit-server.socket

[Service]
# blit server sends sd_notify(READY=1) once the IPC listener is bound,
# so systemctl start blocks until the daemon is actually accepting
# connections — closes the well-known Type=simple race for any unit that
# orders other services After=blit-server.service.
Type=notify
WorkingDirectory=~
ExecStart=blit server
# Keep running after all clients disconnect — PTYs stay alive.
# Stop manually with: systemctl --user stop blit-server.service
Restart=no
# Socket activation passes the fd; BLIT_SOCK is not used.
# Environment inherits from the user session (SHELL, TERM, etc).

[Install]
WantedBy=default.target
