# systemd/blit.service
[Unit]
Description=blit terminal multiplexer
Requires=blit.socket

[Service]
# blit server sends sd_notify(READY=1) once the IPC listener is bound.
Type=notify
WorkingDirectory=~
ExecStart=blit server
Restart=no

[Install]
WantedBy=default.target
